# docker/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
user=root
logfile=/mica/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/mica/logs/nginx.err.log
stdout_logfile=/mica/logs/nginx.out.log

[program:mica]
command=python -u server.py
directory=/mica
autostart=true
autorestart=true
stderr_logfile=/mica/logs/mica.err.log
stdout_logfile=/mica/logs/mica.out.log
environment=PYTHONPATH="/"

[program:mica-ui]
command=python -m mica.demo
directory=/mica
autostart=true
autorestart=true
stderr_logfile=/mica/logs/mica-ui.err.log
stdout_logfile=/mica/logs/mica-ui.out.log
environment=PYTHONPATH="/"